# auto-shutdown/auto-shutdown.service
[Unit]
Description=Check for automatic shutdown conditions

[Service]
ExecStart=/etc/auto-shutdown/auto-shutdown.sh
Type=oneshot
