# auto-shutdown/auto-shutdown.timer
[Unit]
Description=Run automatic shutdown regularly

[Timer]
OnBootSec=10m
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
