# jetson-gpio-shutdown/jetson-gpio-shutdown.service
[Unit]
Description=Automatic Jetson GPIO-triggered shutdown
After=multi-user.target

[Service]
Type=simple
ExecStart=/opt/jetson-gpio-shutdown/jetson-gpio-shutdown.sh
Restart=always
User=root

[Install]
WantedBy=multi-user.target
